# resources/supervisor-seat.ini
[program:seat]
command=:php :artisan queue:work --queue=high,medium,low,default --tries 1 --timeout=86100
process_name = %(program_name)s-80%(process_num)02d
stdout_logfile = /var/log/seat-80%(process_num)02d.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
numprocs=4
directory=:seatdirectory
stopwaitsecs=600
user=:webuser
