# resources/supervisor-seat.ini
[program:seat]
command=:php :artisan horizon
process_name = %(program_name)s
stdout_logfile = /var/log/seat-80%(process_num)02d.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
autostart=true
autorestart=true
numprocs=1
directory=:seatdirectory
stopwaitsecs=600
user=:webuser
